# srcs/requirements/wordpress/conf/supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:php-fpm]
command= sh -c "/usr/sbin/php-fpm7 && kill -s SIGTERM $(cat supervisord.pid)"
autostart=true
autorestart=unexpected
exitcodes=0